# keyboardd/install/keyboardd.service
[Unit]
Description=Daemon for the keyboardd of the MiniBook
ConditionPathExists=/usr/bin/
Requires=sysinit.target
After=sysinit.target basic.target

[Service]

ExecStart=/usr/bin/keyboardd
Restart=no
Type=simple

[Install]
WantedBy=multi-user.target
